# server/systemd/levels-ingest.service
[Unit]
Description=Levels Ingest Job

[Service]
Type=oneshot
WorkingDirectory=/srv/levels
EnvironmentFile=/srv/levels/.env
ExecStart=/usr/bin/python3 /srv/levels/app/ingest.py
