# ops/ingest.timer.tmpl
[Unit]
Description=Run __APP__ ingest every 10 minutes
[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Unit=__APP__-ingest.service
[Install]
WantedBy=timers.target
